# tcz/avahi-2.6.tcz.list
usr/local/bin/avahi-bookmarks
usr/local/bin/avahi-browse
usr/local/bin/avahi-browse-domains
usr/local/bin/avahi-discover-standalone
usr/local/bin/avahi-publish
usr/local/bin/avahi-publish-address
usr/local/bin/avahi-publish-service
usr/local/bin/avahi-resolve
usr/local/bin/avahi-resolve-address
usr/local/bin/avahi-resolve-host-name
usr/local/bin/avahi-set-host-name
usr/local/bin/bshell
usr/local/bin/bssh
usr/local/bin/bvnc
usr/local/etc/avahi/avahi-autoipd.action
usr/local/etc/avahi/avahi-dnsconfd.action
usr/local/etc/avahi/services/sftp-ssh.service
usr/local/etc/avahi/services/ssh.service
usr/local/etc/dbus-1/system.d/avahi-dbus.conf
usr/local/etc/init.d/avahi
usr/local/include/avahi-client/client.h
usr/local/include/avahi-client/lookup.h
usr/local/include/avahi-client/publish.h
usr/local/include/avahi-common/address.h
usr/local/include/avahi-common/alternative.h
usr/local/include/avahi-common/cdecl.h
usr/local/include/avahi-common/defs.h
usr/local/include/avahi-common/domain.h
usr/local/include/avahi-common/error.h
usr/local/include/avahi-common/gccmacro.h
usr/local/include/avahi-common/llist.h
usr/local/include/avahi-common/malloc.h
usr/local/include/avahi-common/rlist.h
usr/local/include/avahi-common/simple-watch.h
usr/local/include/avahi-common/strlst.h
usr/local/include/avahi-common/thread-watch.h
usr/local/include/avahi-common/timeval.h
usr/local/include/avahi-common/watch.h
usr/local/include/avahi-core/core.h
usr/local/include/avahi-core/log.h
usr/local/include/avahi-core/lookup.h
usr/local/include/avahi-core/publish.h
usr/local/include/avahi-core/rr.h
usr/local/include/avahi-glib/glib-malloc.h
usr/local/include/avahi-glib/glib-watch.h
usr/local/include/avahi-gobject/ga-client.h
usr/local/include/avahi-gobject/ga-entry-group.h
usr/local/include/avahi-gobject/ga-enums.h
usr/local/include/avahi-gobject/ga-error.h
usr/local/include/avahi-gobject/ga-record-browser.h
usr/local/include/avahi-gobject/ga-service-browser.h
usr/local/include/avahi-gobject/ga-service-resolver.h
usr/local/include/avahi-ui/avahi-ui.h
usr/local/lib/libavahi-client.a
usr/local/lib/libavahi-client.la
usr/local/lib/libavahi-common.a
usr/local/lib/libavahi-common.la
usr/local/lib/libavahi-core.a
usr/local/lib/libavahi-core.la
usr/local/lib/libavahi-glib.a
usr/local/lib/libavahi-glib.la
usr/local/lib/libavahi-gobject.a
usr/local/lib/libavahi-gobject.la
usr/local/lib/libavahi-ui.a
usr/local/lib/libavahi-ui.la
usr/local/lib/pkgconfig/avahi-client.pc
usr/local/lib/pkgconfig/avahi-core.pc
usr/local/lib/pkgconfig/avahi-glib.pc
usr/local/lib/pkgconfig/avahi-gobject.pc
usr/local/lib/pkgconfig/avahi-ui.pc
usr/local/lib/python2.6/site-packages/avahi/__init__.py
usr/local/lib/python2.6/site-packages/avahi/__init__.pyc
usr/local/lib/python2.6/site-packages/avahi/__init__.pyo
usr/local/sbin/avahi-autoipd
usr/local/sbin/avahi-daemon
usr/local/sbin/avahi-dnsconfd
usr/local/share/applications/avahi-2.6.desktop
usr/local/share/applications/avahi-2.6~1.desktop
usr/local/share/applications/bssh.desktop
usr/local/share/applications/bvnc.desktop
usr/local/share/avahi/avahi-service.dtd
usr/local/share/avahi/files/avahi-daemon.conf
usr/local/share/avahi/files/hosts
usr/local/share/avahi/interfaces/avahi-discover.glade
usr/local/share/avahi/introspection/AddressResolver.introspect
usr/local/share/avahi/introspection/DomainBrowser.introspect
usr/local/share/avahi/introspection/EntryGroup.introspect
usr/local/share/avahi/introspection/HostNameResolver.introspect
usr/local/share/avahi/introspection/RecordBrowser.introspect
usr/local/share/avahi/introspection/Server.introspect
usr/local/share/avahi/introspection/ServiceBrowser.introspect
usr/local/share/avahi/introspection/ServiceResolver.introspect
usr/local/share/avahi/introspection/ServiceTypeBrowser.introspect
usr/local/share/avahi/service-types
usr/local/share/pixmaps/avahi-2.6.png
usr/local/tce.installed/avahi-2.6
