# tcz/bluez.tcz.list
/usr/local/etc/udev/rules.d/97-bluetooth.rules
/usr/local/etc/dbus-1/system.d/bluetooth.conf
/usr/local/share/alsa/bluetooth.conf
/usr/local/sbin/hciemu
/usr/local/sbin/hciconfig
/usr/local/sbin/hciattach
/usr/local/sbin/bluetoothd
/usr/local/bin/dund
/usr/local/bin/pand
/usr/local/bin/hidd
/usr/local/bin/rctest
/usr/local/bin/l2test
/usr/local/bin/ciptool
/usr/local/bin/sdptool
/usr/local/bin/hcitool
/usr/local/bin/l2ping
/usr/local/bin/rfcomm
/usr/local/lib/libbluetooth.so
/usr/local/lib/libbluetooth.so.3
/usr/local/lib/libbluetooth.so.3.10.3
/usr/local/lib/cups/backend/bluetooth
/usr/local/lib/gstreamer-0.10/libgstbluetooth.so
/usr/local/lib/alsa-lib/libasound_module_ctl_bluetooth.so
/usr/local/lib/alsa-lib/libasound_module_pcm_bluetooth.so
/usr/local/bin/simple-agent
/usr/local/bin/simple-service
/usr/local/share/bluez/files/bluetooth/rfcomm.conf
/usr/local/share/bluez/files/bluetooth/main.conf
/usr/local/etc/alsa/example_.asoundrc
/usr/local/etc/init.d/bluez
/usr/local/tce.installed/bluez
