# tcz/cups1311.tcz.list
etc/cups/cupsd.conf
etc/cups/cupsd.conf.N
etc/cups/cupsd.conf.default
etc/cups/mime.convs
etc/cups/mime.types
etc/cups/pstoraster.convs
etc/cups/snmp.conf
etc/cups/snmp.conf.N
etc/dbus-1/system.d/cups.conf
etc/init.d/cups
usr/local/bin/cancel
usr/local/bin/cups-config
usr/local/bin/cupstestdsc
usr/local/bin/cupstestppd
usr/local/bin/lp
usr/local/bin/lpoptions
usr/local/bin/lppasswd
usr/local/bin/lpq
usr/local/bin/lpr
usr/local/bin/lprm
usr/local/bin/lpstat
usr/local/include/cups/adminutil.h
usr/local/include/cups/array.h
usr/local/include/cups/backend.h
usr/local/include/cups/cups.h
usr/local/include/cups/dir.h
usr/local/include/cups/file.h
usr/local/include/cups/http.h
usr/local/include/cups/image.h
usr/local/include/cups/ipp.h
usr/local/include/cups/language.h
usr/local/include/cups/ppd.h
usr/local/include/cups/raster.h
usr/local/include/cups/sidechannel.h
usr/local/include/cups/transcode.h
usr/local/lib/cups/backend/http
usr/local/lib/cups/backend/ipp
usr/local/lib/cups/backend/lpd
usr/local/lib/cups/backend/parallel
usr/local/lib/cups/backend/scsi
usr/local/lib/cups/backend/serial
usr/local/lib/cups/backend/smb
usr/local/lib/cups/backend/snmp
usr/local/lib/cups/backend/socket
usr/local/lib/cups/backend/usb
usr/local/lib/cups/cgi-bin/admin.cgi
usr/local/lib/cups/cgi-bin/classes.cgi
usr/local/lib/cups/cgi-bin/help.cgi
usr/local/lib/cups/cgi-bin/jobs.cgi
usr/local/lib/cups/cgi-bin/printers.cgi
usr/local/lib/cups/daemon/cups-deviced
usr/local/lib/cups/daemon/cups-driverd
usr/local/lib/cups/daemon/cups-lpd
usr/local/lib/cups/daemon/cups-polld
usr/local/lib/cups/filter/gziptoany
usr/local/lib/cups/filter/hpgltops
usr/local/lib/cups/filter/imagetops
usr/local/lib/cups/filter/imagetoraster
usr/local/lib/cups/filter/pdftops
usr/local/lib/cups/filter/pstops
usr/local/lib/cups/filter/rastertodymo
usr/local/lib/cups/filter/rastertoepson
usr/local/lib/cups/filter/rastertohp
usr/local/lib/cups/filter/rastertolabel
usr/local/lib/cups/filter/texttops
usr/local/lib/cups/monitor/bcp
usr/local/lib/cups/monitor/tbcp
usr/local/lib/cups/notifier/mailto
usr/local/lib/cups/notifier/rss
usr/local/lib/libcups.so
usr/local/lib/libcups.so.2
usr/local/lib/libcupsimage.so
usr/local/lib/libcupsimage.so.2
usr/local/sbin/accept
usr/local/sbin/cupsaddsmb
usr/local/sbin/cupsctl
usr/local/sbin/cupsd
usr/local/sbin/cupsdisable
usr/local/sbin/cupsenable
usr/local/sbin/cupsfilter
usr/local/sbin/lpadmin
usr/local/sbin/lpc
usr/local/sbin/lpinfo
usr/local/sbin/lpmove
usr/local/sbin/reject
usr/local/share/cups/banners/classified
usr/local/share/cups/banners/confidential
usr/local/share/cups/banners/secret
usr/local/share/cups/banners/standard
usr/local/share/cups/banners/topsecret
usr/local/share/cups/banners/unclassified
usr/local/share/cups/charmaps/euc-cn.txt
usr/local/share/cups/charmaps/euc-jp.txt
usr/local/share/cups/charmaps/euc-kr.txt
usr/local/share/cups/charmaps/euc-tw.txt
usr/local/share/cups/charmaps/iso-8859-1.txt
usr/local/share/cups/charmaps/iso-8859-10.txt
usr/local/share/cups/charmaps/iso-8859-11.txt
usr/local/share/cups/charmaps/iso-8859-13.txt
usr/local/share/cups/charmaps/iso-8859-14.txt
usr/local/share/cups/charmaps/iso-8859-15.txt
usr/local/share/cups/charmaps/iso-8859-16.txt
usr/local/share/cups/charmaps/iso-8859-2.txt
usr/local/share/cups/charmaps/iso-8859-3.txt
usr/local/share/cups/charmaps/iso-8859-4.txt
usr/local/share/cups/charmaps/iso-8859-5.txt
usr/local/share/cups/charmaps/iso-8859-6.txt
usr/local/share/cups/charmaps/iso-8859-7.txt
usr/local/share/cups/charmaps/iso-8859-8.txt
usr/local/share/cups/charmaps/iso-8859-9.txt
usr/local/share/cups/charmaps/jis-x0213.txt
usr/local/share/cups/charmaps/koi8-r.txt
usr/local/share/cups/charmaps/koi8-u.txt
usr/local/share/cups/charmaps/mac-roman.txt
usr/local/share/cups/charmaps/windows-1250.txt
usr/local/share/cups/charmaps/windows-1251.txt
usr/local/share/cups/charmaps/windows-1252.txt
usr/local/share/cups/charmaps/windows-1253.txt
usr/local/share/cups/charmaps/windows-1254.txt
usr/local/share/cups/charmaps/windows-1255.txt
usr/local/share/cups/charmaps/windows-1256.txt
usr/local/share/cups/charmaps/windows-1257.txt
usr/local/share/cups/charmaps/windows-1258.txt
usr/local/share/cups/charmaps/windows-1361.txt
usr/local/share/cups/charmaps/windows-874.txt
usr/local/share/cups/charmaps/windows-932.txt
usr/local/share/cups/charmaps/windows-936.txt
usr/local/share/cups/charmaps/windows-949.txt
usr/local/share/cups/charmaps/windows-950.txt
usr/local/share/cups/charsets/utf-8
usr/local/share/cups/data/HPGLprolog
usr/local/share/cups/data/psglyphs
usr/local/share/cups/data/testprint.ps
usr/local/share/cups/fonts/Courier
usr/local/share/cups/fonts/Courier-Bold
usr/local/share/cups/fonts/Courier-BoldOblique
usr/local/share/cups/fonts/Courier-Oblique
usr/local/share/cups/fonts/Symbol
usr/local/share/cups/model/deskjet.ppd
usr/local/share/cups/model/dymo.ppd
usr/local/share/cups/model/epson24.ppd
usr/local/share/cups/model/epson9.ppd
usr/local/share/cups/model/intelbar.ppd
usr/local/share/cups/model/laserjet.ppd
usr/local/share/cups/model/okidat24.ppd
usr/local/share/cups/model/okidata9.ppd
usr/local/share/cups/model/stcolor.ppd
usr/local/share/cups/model/stcolor2.ppd
usr/local/share/cups/model/stphoto.ppd
usr/local/share/cups/model/stphoto2.ppd
usr/local/share/cups/model/zebra.ppd
usr/local/share/cups/model/zebracpl.ppd
usr/local/share/cups/model/zebraep1.ppd
usr/local/share/cups/model/zebraep2.ppd
usr/local/share/cups/templates/add-class.tmpl
usr/local/share/cups/templates/add-printer.tmpl
usr/local/share/cups/templates/add-rss-subscription.tmpl
usr/local/share/cups/templates/admin.tmpl
usr/local/share/cups/templates/choose-device.tmpl
usr/local/share/cups/templates/choose-make.tmpl
usr/local/share/cups/templates/choose-model.tmpl
usr/local/share/cups/templates/choose-serial.tmpl
usr/local/share/cups/templates/choose-uri.tmpl
usr/local/share/cups/templates/class-added.tmpl
usr/local/share/cups/templates/class-confirm.tmpl
usr/local/share/cups/templates/class-deleted.tmpl
usr/local/share/cups/templates/class-jobs-header.tmpl
usr/local/share/cups/templates/class-modified.tmpl
usr/local/share/cups/templates/classes-header.tmpl
usr/local/share/cups/templates/classes.tmpl
usr/local/share/cups/templates/edit-config.tmpl
usr/local/share/cups/templates/error-op.tmpl
usr/local/share/cups/templates/error.tmpl
usr/local/share/cups/templates/header.tmpl
usr/local/share/cups/templates/help-header.tmpl
usr/local/share/cups/templates/help-printable.tmpl
usr/local/share/cups/templates/job-cancel.tmpl
usr/local/share/cups/templates/job-hold.tmpl
usr/local/share/cups/templates/job-move.tmpl
usr/local/share/cups/templates/job-moved.tmpl
usr/local/share/cups/templates/job-release.tmpl
usr/local/share/cups/templates/job-restart.tmpl
usr/local/share/cups/templates/jobs-header.tmpl
usr/local/share/cups/templates/jobs.tmpl
usr/local/share/cups/templates/list-available-printers.tmpl
usr/local/share/cups/templates/maintenance.tmpl
usr/local/share/cups/templates/modify-class.tmpl
usr/local/share/cups/templates/modify-printer.tmpl
usr/local/share/cups/templates/norestart.tmpl
usr/local/share/cups/templates/option-boolean.tmpl
usr/local/share/cups/templates/option-conflict.tmpl
usr/local/share/cups/templates/option-header.tmpl
usr/local/share/cups/templates/option-pickmany.tmpl
usr/local/share/cups/templates/option-pickone.tmpl
usr/local/share/cups/templates/option-trailer.tmpl
usr/local/share/cups/templates/pager.tmpl
usr/local/share/cups/templates/printer-accept.tmpl
usr/local/share/cups/templates/printer-added.tmpl
usr/local/share/cups/templates/printer-configured.tmpl
usr/local/share/cups/templates/printer-confirm.tmpl
usr/local/share/cups/templates/printer-default.tmpl
usr/local/share/cups/templates/printer-deleted.tmpl
usr/local/share/cups/templates/printer-jobs-header.tmpl
usr/local/share/cups/templates/printer-modified.tmpl
usr/local/share/cups/templates/printer-purge.tmpl
usr/local/share/cups/templates/printer-reject.tmpl
usr/local/share/cups/templates/printer-start.tmpl
usr/local/share/cups/templates/printer-stop.tmpl
usr/local/share/cups/templates/printers-header.tmpl
usr/local/share/cups/templates/printers.tmpl
usr/local/share/cups/templates/restart.tmpl
usr/local/share/cups/templates/samba-export.tmpl
usr/local/share/cups/templates/samba-exported.tmpl
usr/local/share/cups/templates/search.tmpl
usr/local/share/cups/templates/set-printer-options-header.tmpl
usr/local/share/cups/templates/set-printer-options-trailer.tmpl
usr/local/share/cups/templates/subscription-added.tmpl
usr/local/share/cups/templates/subscription-canceled.tmpl
usr/local/share/cups/templates/test-page.tmpl
usr/local/share/cups/templates/trailer.tmpl
usr/local/share/cups/templates/users.tmpl
usr/local/share/doc/cups/cups-printable.css
usr/local/share/doc/cups/cups.css
usr/local/share/doc/cups/favicon.ico
usr/local/share/doc/cups/help/accounting.html
usr/local/share/doc/cups/help/api-array.html
usr/local/share/doc/cups/help/api-cups.html
usr/local/share/doc/cups/help/api-filedir.html
usr/local/share/doc/cups/help/api-filter.html
usr/local/share/doc/cups/help/api-httpipp.html
usr/local/share/doc/cups/help/api-ppd.html
usr/local/share/doc/cups/help/api-raster.html
usr/local/share/doc/cups/help/cgi.html
usr/local/share/doc/cups/help/glossary.html
usr/local/share/doc/cups/help/kerberos.html
usr/local/share/doc/cups/help/license.html
usr/local/share/doc/cups/help/man-accept.html
usr/local/share/doc/cups/help/man-backend.html
usr/local/share/doc/cups/help/man-cancel.html
usr/local/share/doc/cups/help/man-classes.conf.html
usr/local/share/doc/cups/help/man-client.conf.html
usr/local/share/doc/cups/help/man-cups-config.html
usr/local/share/doc/cups/help/man-cups-lpd.html
usr/local/share/doc/cups/help/man-cups-polld.html
usr/local/share/doc/cups/help/man-cupsaddsmb.html
usr/local/share/doc/cups/help/man-cupsd.conf.html
usr/local/share/doc/cups/help/man-cupsd.html
usr/local/share/doc/cups/help/man-cupsenable.html
usr/local/share/doc/cups/help/man-cupstestdsc.html
usr/local/share/doc/cups/help/man-cupstestppd.html
usr/local/share/doc/cups/help/man-filter.html
usr/local/share/doc/cups/help/man-lp.html
usr/local/share/doc/cups/help/man-lpadmin.html
usr/local/share/doc/cups/help/man-lpc.html
usr/local/share/doc/cups/help/man-lpinfo.html
usr/local/share/doc/cups/help/man-lpmove.html
usr/local/share/doc/cups/help/man-lpoptions.html
usr/local/share/doc/cups/help/man-lppasswd.html
usr/local/share/doc/cups/help/man-lpq.html
usr/local/share/doc/cups/help/man-lpr.html
usr/local/share/doc/cups/help/man-lprm.html
usr/local/share/doc/cups/help/man-lpstat.html
usr/local/share/doc/cups/help/man-mime.convs.html
usr/local/share/doc/cups/help/man-mime.types.html
usr/local/share/doc/cups/help/man-printers.conf.html
usr/local/share/doc/cups/help/network.html
usr/local/share/doc/cups/help/options.html
usr/local/share/doc/cups/help/overview.html
usr/local/share/doc/cups/help/policies.html
usr/local/share/doc/cups/help/ref-access_log.html
usr/local/share/doc/cups/help/ref-classes-conf.html
usr/local/share/doc/cups/help/ref-client-conf.html
usr/local/share/doc/cups/help/ref-cupsd-conf.html
usr/local/share/doc/cups/help/ref-error_log.html
usr/local/share/doc/cups/help/ref-mailto-conf.html
usr/local/share/doc/cups/help/ref-page_log.html
usr/local/share/doc/cups/help/ref-printers-conf.html
usr/local/share/doc/cups/help/ref-snmp-conf.html
usr/local/share/doc/cups/help/ref-subscriptions-conf.html
usr/local/share/doc/cups/help/security.html
usr/local/share/doc/cups/help/spec-cmp.html
usr/local/share/doc/cups/help/spec-command.html
usr/local/share/doc/cups/help/spec-design.html
usr/local/share/doc/cups/help/spec-ipp.html
usr/local/share/doc/cups/help/spec-postscript.html
usr/local/share/doc/cups/help/spec-ppd.html
usr/local/share/doc/cups/help/spec-raster.html
usr/local/share/doc/cups/help/spec-stp.html
usr/local/share/doc/cups/help/standard.html
usr/local/share/doc/cups/help/translation.html
usr/local/share/doc/cups/help/whatsnew.html
usr/local/share/doc/cups/images/bottom-left.gif
usr/local/share/doc/cups/images/bottom-right.gif
usr/local/share/doc/cups/images/button-accept-jobs.gif
usr/local/share/doc/cups/images/button-add-class.gif
usr/local/share/doc/cups/images/button-add-printer.gif
usr/local/share/doc/cups/images/button-add-rss-subscription.gif
usr/local/share/doc/cups/images/button-add-this-printer.gif
usr/local/share/doc/cups/images/button-cancel-all-jobs.gif
usr/local/share/doc/cups/images/button-cancel-job.gif
usr/local/share/doc/cups/images/button-cancel-subscription.gif
usr/local/share/doc/cups/images/button-change-settings.gif
usr/local/share/doc/cups/images/button-clean-print-heads.gif
usr/local/share/doc/cups/images/button-clear.gif
usr/local/share/doc/cups/images/button-continue.gif
usr/local/share/doc/cups/images/button-delete-class.gif
usr/local/share/doc/cups/images/button-delete-printer.gif
usr/local/share/doc/cups/images/button-edit-configuration-file.gif
usr/local/share/doc/cups/images/button-export-samba.gif
usr/local/share/doc/cups/images/button-find-new-printers.gif
usr/local/share/doc/cups/images/button-help.gif
usr/local/share/doc/cups/images/button-hold-job.gif
usr/local/share/doc/cups/images/button-manage-classes.gif
usr/local/share/doc/cups/images/button-manage-jobs.gif
usr/local/share/doc/cups/images/button-manage-printers.gif
usr/local/share/doc/cups/images/button-manage-server.gif
usr/local/share/doc/cups/images/button-modify-class.gif
usr/local/share/doc/cups/images/button-modify-printer.gif
usr/local/share/doc/cups/images/button-move-job.gif
usr/local/share/doc/cups/images/button-move-jobs.gif
usr/local/share/doc/cups/images/button-print-self-test-page.gif
usr/local/share/doc/cups/images/button-print-test-page.gif
usr/local/share/doc/cups/images/button-publish-printer.gif
usr/local/share/doc/cups/images/button-reject-jobs.gif
usr/local/share/doc/cups/images/button-release-job.gif
usr/local/share/doc/cups/images/button-restart-job.gif
usr/local/share/doc/cups/images/button-save-changes.gif
usr/local/share/doc/cups/images/button-search.gif
usr/local/share/doc/cups/images/button-set-allowed-users.gif
usr/local/share/doc/cups/images/button-set-as-default.gif
usr/local/share/doc/cups/images/button-set-printer-options.gif
usr/local/share/doc/cups/images/button-show-active.gif
usr/local/share/doc/cups/images/button-show-all.gif
usr/local/share/doc/cups/images/button-show-completed.gif
usr/local/share/doc/cups/images/button-show-next.gif
usr/local/share/doc/cups/images/button-show-previous.gif
usr/local/share/doc/cups/images/button-sort-ascending.gif
usr/local/share/doc/cups/images/button-sort-descending.gif
usr/local/share/doc/cups/images/button-start-class.gif
usr/local/share/doc/cups/images/button-start-printer.gif
usr/local/share/doc/cups/images/button-stop-class.gif
usr/local/share/doc/cups/images/button-stop-printer.gif
usr/local/share/doc/cups/images/button-unpublish-printer.gif
usr/local/share/doc/cups/images/button-use-default-config.gif
usr/local/share/doc/cups/images/button-view-access-log.gif
usr/local/share/doc/cups/images/button-view-error-log.gif
usr/local/share/doc/cups/images/button-view-page-log.gif
usr/local/share/doc/cups/images/button-view-printable-version.gif
usr/local/share/doc/cups/images/classes.gif
usr/local/share/doc/cups/images/cups-block-diagram.gif
usr/local/share/doc/cups/images/esp-logo.gif
usr/local/share/doc/cups/images/happy.gif
usr/local/share/doc/cups/images/logo.gif
usr/local/share/doc/cups/images/printer-idle.gif
usr/local/share/doc/cups/images/printer-processing.gif
usr/local/share/doc/cups/images/printer-stopped.gif
usr/local/share/doc/cups/images/raster.png
usr/local/share/doc/cups/images/raster.svg
usr/local/share/doc/cups/images/smiley.jpg
usr/local/share/doc/cups/images/tab-left.gif
usr/local/share/doc/cups/images/tab-right.gif
usr/local/share/doc/cups/images/top-left.gif
usr/local/share/doc/cups/images/top-middle.gif
usr/local/share/doc/cups/images/top-right.gif
usr/local/share/doc/cups/images/web-add-this-printer.gif
usr/local/share/doc/cups/images/web-password.gif
usr/local/share/doc/cups/images/web-printer-driver.gif
usr/local/share/doc/cups/images/web-printer-status.gif
usr/local/share/doc/cups/images/web-set-printer-options.gif
usr/local/share/doc/cups/images/web-test-page.gif
usr/local/share/doc/cups/index.html
usr/local/share/doc/cups/robots.txt
