# tcz/open-vm-tools-x.tcz.list
usr/local/bin/vmtoolsd
usr/local/bin/vmware-checkvm
usr/local/bin/vmware-hgfsclient
usr/local/bin/vmware-rpctool
usr/local/bin/vmware-toolbox
usr/local/bin/vmware-toolbox-cmd
usr/local/bin/vmware-user-suid-wrapper
usr/local/bin/vmware-vmblock-fuse
usr/local/bin/vmware-xferlogs
usr/local/etc/init.d/open-vm-tools-x
usr/local/etc/pam.d/vmtoolsd
usr/local/etc/vmware-tools/plugins
usr/local/etc/vmware-tools/poweroff-vm-default
usr/local/etc/vmware-tools/poweron-vm-default
usr/local/etc/vmware-tools/resume-vm-default
usr/local/etc/vmware-tools/scripts/vmware/network
usr/local/etc/vmware-tools/statechange.subr
usr/local/etc/vmware-tools/suspend-vm-default
usr/local/etc/vmware-tools/vm-support
usr/local/etc/xdg/autostart/vmware-user.desktop
usr/local/lib/libguestlib.so
usr/local/lib/libguestlib.so.0
usr/local/lib/libguestlib.so.0.0.0
usr/local/lib/libhgfs.so
usr/local/lib/libhgfs.so.0
usr/local/lib/libhgfs.so.0.0.0
usr/local/lib/libvmtools.so
usr/local/lib/libvmtools.so.0
usr/local/lib/libvmtools.so.0.0.0
usr/local/lib/open-vm-tools/plugins/common/libhgfsServer.so
usr/local/lib/open-vm-tools/plugins/common/libvix.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libguestInfo.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libpowerOps.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libtimeSync.so
usr/local/lib/open-vm-tools/plugins/vmsvc/libvmbackup.so
usr/local/lib/open-vm-tools/plugins/vmusr/libdesktopEvents.so
usr/local/lib/open-vm-tools/plugins/vmusr/libdndcp.so
usr/local/lib/open-vm-tools/plugins/vmusr/libresolutionSet.so
usr/local/lib/open-vm-tools/plugins/vmusr/libunity.so
usr/local/sbin/mount.vmhgfs
usr/local/share/open-vm-tools/messages/de/toolboxcmd.vmsg
usr/local/share/open-vm-tools/messages/de/vmtoolsd.vmsg
usr/local/share/open-vm-tools/messages/ja/toolboxcmd.vmsg
usr/local/share/open-vm-tools/messages/ja/vmtoolsd.vmsg
usr/local/share/open-vm-tools/messages/ko/toolboxcmd.vmsg
usr/local/share/open-vm-tools/messages/ko/vmtoolsd.vmsg
usr/local/share/open-vm-tools/messages/zh_CN/toolboxcmd.vmsg
usr/local/tce.installed/open-vm-tools-x
